# 800-how-to-run-on-production/service.supervisor.ini
[program:gwservice]
command=/home/exu/src/github.com/exu/go-workshops/todo-800-services/service
